# buildroot-2012.02/target/device/geep/dban/skeleton/etc/udev/rules.d/z60_firmware.rules
# firmware
ACTION=="add", SUBSYSTEM=="firmware", ENV{FIRMWARE}=="*", RUN+="/lib/udev/firmware.sh"
